# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:auth-service]
directory=/app/backend/services/auth-service
command=node src/server.js
priority=10
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",JWT_SECRET="epfms-local-JWT-MUST-MATCH-ALL-SERVICES-2026",JWT_EXPIRES_IN="8h",PORT="3001",DATABASE_PATH="/data/auth.db"

[program:user-service]
directory=/app/backend/services/user-service
command=node src/server.js
priority=15
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",INTERNAL_SERVICE_KEY="epfms-local-INTERNAL-KEY-MATCH-2026",PORT="3002",DATABASE_PATH="/data/user.db",NOTIFICATION_SERVICE_URL="http://127.0.0.1:3003"

[program:notification-service]
directory=/app/backend/services/notification-service
command=node src/server.js
priority=20
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",JWT_SECRET="epfms-local-JWT-MUST-MATCH-ALL-SERVICES-2026",JWT_EXPIRES_IN="8h",INTERNAL_SERVICE_KEY="epfms-local-INTERNAL-KEY-MATCH-2026",PORT="3003",DATABASE_PATH="/data/notification.db"

[program:performance-review-service]
directory=/app/backend/services/performance-review-service
command=node src/server.js
priority=25
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",JWT_SECRET="epfms-local-JWT-MUST-MATCH-ALL-SERVICES-2026",JWT_EXPIRES_IN="8h",INTERNAL_SERVICE_KEY="epfms-local-INTERNAL-KEY-MATCH-2026",PORT="8001",DATABASE_PATH="/data/performance.db",NOTIFICATION_SERVICE_URL="http://127.0.0.1:3003",USER_SERVICE_URL="http://127.0.0.1:3002"

[program:feedback-service]
directory=/app/backend/services/feedback-service
command=node src/server.js
priority=30
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",JWT_SECRET="epfms-local-JWT-MUST-MATCH-ALL-SERVICES-2026",JWT_EXPIRES_IN="8h",INTERNAL_SERVICE_KEY="epfms-local-INTERNAL-KEY-MATCH-2026",PORT="8002",DATABASE_PATH="/data/feedback.db",USER_SERVICE_URL="http://127.0.0.1:3002",NOTIFICATION_SERVICE_URL="http://127.0.0.1:3003",ANALYTICS_SERVICE_URL="http://127.0.0.1:8004"

[program:analytics-service]
directory=/app/backend/services/analytics-service
command=node src/server.js
priority=35
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",JWT_SECRET="epfms-local-JWT-MUST-MATCH-ALL-SERVICES-2026",JWT_EXPIRES_IN="8h",INTERNAL_SERVICE_KEY="epfms-local-INTERNAL-KEY-MATCH-2026",PORT="8004",DATABASE_PATH="/data/analytics.db"

[program:ai-insights-service]
directory=/app/backend/services/ai-insights-service
command=node src/server.js
priority=40
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",PORT="8005",AI_SERVICE_KEY="epfms-local-ai-key-change-me"

[program:api-gateway]
directory=/app/backend/services/api-gateway
command=node src/server.js
priority=50
startsecs=3
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",LOG_LEVEL="info",JWT_SECRET="epfms-local-JWT-MUST-MATCH-ALL-SERVICES-2026",JWT_EXPIRES_IN="8h",INTERNAL_SERVICE_KEY="epfms-local-INTERNAL-KEY-MATCH-2026",GATEWAY_PORT="8080",AUTH_SERVICE_URL="http://127.0.0.1:3001",USER_SERVICE_URL="http://127.0.0.1:3002",NOTIFICATION_SERVICE_URL="http://127.0.0.1:3003",PERFORMANCE_SERVICE_URL="http://127.0.0.1:8001",FEEDBACK_SERVICE_URL="http://127.0.0.1:8002",ANALYTICS_SERVICE_URL="http://127.0.0.1:8004",AI_INSIGHTS_SERVICE_URL="http://127.0.0.1:8005",AI_SERVICE_KEY="epfms-local-ai-key-change-me",RATE_LIMIT_WINDOW_MS="60000",RATE_LIMIT_MAX="200"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=60
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
